# debian/droidian-quirks-device.droidian-boot-wlan.path
[Path]
PathChanged=/dev/wmtWifi

[Install]
# FIXME: multi-user.target is not the right target for this, but
# we're bound on the lxc container's target.
WantedBy=multi-user.target
