# quirks-33/systemd/android-mount.service.d/10-selinux.conf
[Service]
ExecStartPost=/usr/bin/mount -o bind /android/system/lib64/libselinux_stubs.so /var/lib/lxc/android/rootfs/vendor/lib64/libselinux_stubs.so
