# keyring/systemd/user/gnome-keyring-daemon.service.d/90-initial-setup.conf
[Unit]
# Do not start keyring daemon if initial setup hasn't run
ConditionPathExists=%E/gnome-initial-setup-done
