# keyring/systemd/system/phosh.service.d/90-pamname.conf
[Service]
# Use phosh specific configuration
PAMName=phosh
