# initial-setup/systemd/user/gnome-initial-setup-first-login.service.d/10-droidian.conf
[Service]
Environment=LANG=en_US.UTF8
ExecStart=
ExecStart=/usr/libexec/gnome-initial-setup
#ExecStart=/usr/bin/pkill -SIGTERM gnome-session
