# initial-setup/systemd/user/gnome-initial-setup-first-login.service.d/10-droidian.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/gnome-initial-setup
